# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:athens-server]
command=node --max-old-space-size=8192 index.js
directory=/app/Athens-server
environment=PORT=8979,HOST="0.0.0.0",MONGO_URL="%(ENV_MONGO_URL)s",MONGO_DB="%(ENV_MONGO_DB)s",API_KEYS_ENCRYPTION_KEY="%(ENV_API_KEYS_ENCRYPTION_KEY)s",AI_BFF_URL="http://127.0.0.1:3920",PUPPETEER_CACHE_DIR="%(ENV_PUPPETEER_CACHE_DIR)s",PUPPETEER_ARGS="%(ENV_PUPPETEER_ARGS)s",JSON_BODY_LIMIT="4096mb"
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:avalon-relay]
command=node src/server.js
directory=/app/project-avalon/packages/backend
environment=AVALON_PORT=3847,HOST="0.0.0.0",CORS_ORIGIN="*",AVALON_PING_INTERVAL_MS=25000,AVALON_PING_TIMEOUT_MS=60000
autostart=true
autorestart=true
priority=25
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:ai-bff]
command=node --max-old-space-size=4096 dist/index.js
directory=/app/ai-bff
environment=PORT=3920,MONGO_URL="%(ENV_MONGO_URL)s",MONGO_DB="%(ENV_MONGO_DB)s",CORS_ORIGIN="*"
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=100
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
